# local-modules/desktop/templates/x11/startx.service.erb
[Unit]
Description=X11 session for <%= @user %>
After=graphical.target systemd-user-sessions.service

[Service]
Type=simple
User=<%= @user %>
WorkingDirectory=~
PAMName=login
StandardInput=tty
TTYPath=/dev/tty8
UtmpIdentifier=tty8
ExecStartPre=/usr/bin/chvt 8
ExecStart=/usr/bin/startx
UnsetEnvironment=TERM
Environment=XDG_SESSION_TYPE=x11

[Install]
WantedBy=graphical.target
